# user_scripts/waybar/network/network_meter.service
[Unit]
Description=Waybar network speed daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Use %h for home directory
ExecStart=%h/user_scripts/waybar/network/network_meter_daemon.sh
Restart=always
RestartSec=1s

# Hardening (Optional but recommended for always-on scripts)
ProtectSystem=full
ProtectHome=read-only
# Allow writing to XDG_RUNTIME_DIR (implied by default in user services usually, but good to be safe)
ReadWritePaths=%t

[Install]
WantedBy=default.target
