# user_scripts/battery/notify/battery_notify.service
[Unit]
Description=Battery Status Notification Service
# Ensure we start AFTER the graphical session is truly ready
After=graphical-session.target
# If you use UWSM, this binds the service to the UWSM session
PartOf=graphical-session.target

[Service]
Type=simple
# Use %h for home directory compatibility
ExecStart=%h/user_scripts/battery/notify/battery_notify.sh
Restart=on-failure
RestartSec=5

# IMPORTANT: Do not hardcode DISPLAY or DBUS. 
# Systemd user session should inherit these from UWSM's activation.
# However, to be safe, we enable logging:
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
